# root/etc/init.d/prefix-nat
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=18
STOP=99

CONFIGURATION="prefix-nat"
NFT_FILE="/usr/share/nftables.d/table-post/10-prefix-nat.nft"

start_service() {
	config_load "${CONFIGURATION}"

	local global_enabled
	config_get_bool global_enabled global enabled 0
	if [ "$global_enabled" -eq 0 ]; then
		rm -f "$NFT_FILE"
		fw4 reload 2>/dev/null
		return 0
	fi

	generate_nft
	fw4 reload 2>/dev/null
}

stop_service() {
	rm -f "$NFT_FILE"
	fw4 reload 2>/dev/null
}

generate_nft() {
	local has_pre=0 has_post=0

	mkdir -p "$(dirname "$NFT_FILE")"

	# Collect rules first to check if we have any
	local pre_rules="" post_rules=""
	config_foreach collect_rule rule

	if [ -z "$pre_rules" ] && [ -z "$post_rules" ]; then
		rm -f "$NFT_FILE"
		return 0
	fi

	{
		echo "# Auto-generated by /etc/init.d/prefix-nat"
		echo "# Do not edit — managed via UCI (/etc/config/prefix-nat)"
		echo ""

		if [ -n "$pre_rules" ]; then
			echo "chain prefix_nat_prerouting {"
			echo "	type nat hook prerouting priority dstnat - 1; policy accept;"
			echo "$pre_rules"
			echo "}"
			echo ""
		fi

		if [ -n "$post_rules" ]; then
			echo "chain prefix_nat_postrouting {"
			echo "	type nat hook postrouting priority srcnat + 1; policy accept;"
			echo "$post_rules"
			echo "}"
		fi
	} > "${NFT_FILE}.tmp"
	mv "${NFT_FILE}.tmp" "$NFT_FILE"

	logger -t prefix-nat "Generated $NFT_FILE"
}

collect_rule() {
	local section="$1"
	local enabled src dest interface remote description

	config_get_bool enabled "$section" enabled 1
	[ "$enabled" -eq 0 ] && return 0

	config_get src "$section" src
	config_get dest "$section" dest
	config_get interface "$section" interface
	config_get remote "$section" remote ""
	config_get description "$section" description ""

	[ -z "$src" ] || [ -z "$dest" ] || [ -z "$interface" ] && {
		logger -t prefix-nat "Incomplete rule '$section': src='$src' dest='$dest' interface='$interface'"
		return 0
	}

	# Validate CIDR notation (octets 0-255, prefix 0-32)
	validate_cidr() {
		echo "$1" | grep -qE '^((25[0-5]|2[0-4][0-9]|[01]?[0-9][0-9]?)\.){3}(25[0-5]|2[0-4][0-9]|[01]?[0-9][0-9]?)/(3[0-2]|[12]?[0-9])$'
	}
	validate_cidr "$src" || {
		logger -t prefix-nat "Invalid src subnet '$src' in rule '$section'"
		return 0
	}
	validate_cidr "$dest" || {
		logger -t prefix-nat "Invalid dest subnet '$dest' in rule '$section'"
		return 0
	}
	if [ -n "$remote" ]; then
		validate_cidr "$remote" || {
			logger -t prefix-nat "Invalid remote subnet '$remote' in rule '$section'"
			return 0
		}
	fi

	# Validate prefix lengths match (required by nftables prefix NAT)
	local src_prefix="${src##*/}"
	local dest_prefix="${dest##*/}"
	[ "$src_prefix" != "$dest_prefix" ] && {
		logger -t prefix-nat "Prefix length mismatch in rule '$section': src=/$src_prefix dest=/$dest_prefix"
		return 0
	}

	# Validate interface name (alphanumeric, underscore, hyphen, dot only)
	echo "$interface" | grep -qE '^[a-zA-Z0-9_.-]+$' || {
		logger -t prefix-nat "Invalid interface name '$interface' in rule '$section'"
		return 0
	}

	# Sanitize description for safe nftables embedding
	local comment=""
	if [ -n "$description" ]; then
		description=$(printf '%s' "$description" | tr -d '";{}\\' | head -c 64)
		comment=" comment \"$description\""
	fi

	# Build optional remote subnet match
	local pre_remote="" post_remote=""
	if [ -n "$remote" ]; then
		pre_remote=" ip saddr ${remote}"
		post_remote=" ip daddr ${remote}"
	fi

	# PREROUTING: inbound traffic on this interface, [from remote,] dst=dest -> translate to src
	pre_rules="${pre_rules}
	iifname \"${interface}\"${pre_remote} ip daddr ${dest} dnat ip prefix to ${src}${comment}"

	# POSTROUTING: outbound traffic on this interface, src=src [to remote] -> translate to dest
	post_rules="${post_rules}
	oifname \"${interface}\" ip saddr ${src}${post_remote} snat ip prefix to ${dest}${comment}"
}

service_triggers() {
	procd_add_reload_trigger "${CONFIGURATION}"
}

reload_service() {
	start_service
}
